# scripts/com.claudeskills.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeskills.scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/erichowens/coding/some_claude_skills/scripts/scan-and-ingest-skills.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/erichowens/coding/some_claude_skills/logs/skill-scan-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/erichowens/coding/some_claude_skills/logs/skill-scan-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/erichowens/coding/some_claude_skills</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
